# recipes-apps/hello_cpp/files/hello.service
[Unit]
Description=GNU Hello World startup script

[Service]
ExecStart=/usr/sbin/hellocpp

[Install]
WantedBy=multi-user.target